# rt/config/rt-update-index.service
[Unit]
After=httpd

[Service]
Type=oneshot
ExecStart=/opt/rt4/sbin/rt-fulltext-indexer
User=apache
Group=apache
